# cron/crontab
# m h dom mon dow user  command

3 */2 * * *	loxberry /usr/bin/php REPLACELBPHTMLDIR/system/battery.php > /dev/null 2>&1
