# webfrontend/html/src/Core/systemd/sonos_event_listener.service
[Unit]
Description=LoxBerry Sonos Event Listener (MQTT)
StartLimitIntervalSec=60
StartLimitBurst=10
After=network-online.target mosquitto.service
Wants=network-online.target
ConditionPathExists=/opt/loxberry/webfrontend/html/plugins/sonos4lox/src/Core/Event/EventHandler.php

[Service]
Type=simple
User=loxberry
Group=loxberry
WorkingDirectory=/opt/loxberry/webfrontend/html/plugins/sonos4lox
ExecStart=/usr/bin/php /opt/loxberry/webfrontend/html/plugins/sonos4lox/src/Core/Event/EventHandler.php
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
